# deploy/compose/systemd/base-validator-image-updater@.service
[Unit]
Description=Base validator image auto-update once for project %i
Documentation=file:///usr/local/lib/base/validator-image-updater.sh
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
# Instance name (%i) is the Compose project name (e.g. base-validator-5gzi).
Environment=COMPOSE_PROJECT_NAME=%i
# Paths + track/hold only. Do NOT EnvironmentFile the project .env: compose
# prefers process-env IMAGE pins over --env-file, which freezes auto-update
# to the pre-tick digest and never applies the newly rewritten pin.
EnvironmentFile=-/etc/base/validator-image-updater/%i.env
ExecStart=/usr/local/lib/base/validator-image-updater.sh once
Nice=10
TimeoutStartSec=600
# Host-only helper. The long-lived agent container never gets docker.sock.

[Install]
WantedBy=multi-user.target
